# etc/sv/nginx/run
#!/bin/sh
exec 2>&1

exec chpst -P /usr/sbin/nginx -p /var/lib/nginx
